# deploy/lightpn-hub.service
[Unit]
Description=LightPN hub (control plane + admin panel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/lightpn-hub --data-dir /var/lib/lightpn/hub
Restart=always
RestartSec=3

# Hardening. Runs as root so the data dir is shared cleanly with the
# `lightpn-hub admin set-password` CLI (run via sudo); ProtectSystem
# still confines writes to the state directory.
StateDirectory=lightpn/hub
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
